# net-misc/geomyidae/files/geomyidae.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
command="geomyidae"
command_args="-o 70 -p 7000"

depend() {
	need net
}
